# systemd/secrets-audit-email.service
[Unit]
Description=Email Jeff a monthly secrets-rotation PDF audit
After=docker.service network-online.target
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/dev-ops
ExecStart=/opt/dev-ops/scripts/cron/secrets-audit-email.sh
TimeoutStartSec=180
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
