# systemd/secrets-audit-email.timer
[Unit]
Description=Monthly secrets-rotation PDF audit (1st of month, 09:30 UTC)

[Timer]
OnCalendar=*-*-01 09:30 UTC
Persistent=true

[Install]
WantedBy=timers.target
