# systemd/smoke-report-email.timer
[Unit]
Description=Weekly smoke-test PDF digest (Sundays 18:00 UTC)

[Timer]
OnCalendar=Sun 18:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
