# netcup/gitea/watchdog-stale-runs.timer
[Unit]
Description=Run gitea-actions stale-Running watchdog every 15 min
Documentation=file:///opt/scripts/watchdog-stale-runs.sh

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
AccuracySec=30s
Persistent=true

[Install]
WantedBy=timers.target
