# netcup/scripts/uptime-kuma-tier0-probe.timer
[Unit]
Description=Run Tier 0 saturation probe every 5 minutes

[Timer]
OnBootSec=120s
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
