# netcup/uptime-kuma/uptime-kuma-engine-pool-probe.service
[Unit]
Description=Uptime Kuma morpheus-engine-pool latency probe
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/opt/scripts/uptime-kuma-engine-pool-probe.sh
# Per-engine prior-failed-count state lives here; auto-created with
# the right ownership when systemd starts the unit.
StateDirectory=uptime-kuma-engine-pool-probe
# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
