# netcup/uptime-kuma/uptime-kuma-host-probe.service
[Unit]
Description=Uptime Kuma host health probe
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/opt/scripts/uptime-kuma-host-probe.sh
StateDirectory=uptime-kuma-host-probe
# Bound runtime so a wedged invocation can't block the next tick
TimeoutStartSec=60
# Get scheduled even when the host is thrashing on CPU/IO
Nice=-5
IOSchedulingClass=best-effort
IOSchedulingPriority=2
# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
