# netcup/uptime-kuma/uptime-kuma-host-probe.timer
[Unit]
Description=Run Uptime Kuma host health probe every 2 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=2min
AccuracySec=10s
# Catch up if the host was under load and a tick was missed
Persistent=true

[Install]
WantedBy=timers.target
